# luzern-velox-vibebot.service
[Unit]
Description=Luzern Velox Vibebot Telegram Bot
# Start after the network is available
After=network.target

[Service]
# User and Group that the service will run as
# Ensure this user has read/write permissions for WorkingDirectory and the executable
User=your_user
Group=your_group

# Set the working directory to the project root
WorkingDirectory=/home/alessandro/git/luzern-velox-vibebot # <-- ADJUST THIS PATH

# Environment variables (TELEGRAM_BOT_TOKEN, TELEGRAM_CHAT_ID, RUST_LOG)
# are expected to be loaded from the .env file located in the WorkingDirectory.
# Ensure the .env file exists and is configured correctly.
# Alternatively, uncomment and use EnvironmentFile= to specify a different path,
# or use Environment="VAR=value" directives here (less recommended for secrets).
# Example: EnvironmentFile=/etc/luzern-velox-vibebot/config
# Example: Environment="RUST_LOG=debug"

# Command to execute
# Ensure the binary is built and located here
ExecStart=/home/alessandro/git/luzern-velox-vibebot/target/release/luzern-velox-vibebot # <-- ADJUST THIS PATH if needed

# Restart the service if it fails
Restart=on-failure
RestartSec=5s

# Standard output and error logging configuration
StandardOutput=journal
StandardError=journal

[Install]
# Enable the service for the default multi-user target
WantedBy=multi-user.target
